# N1/files/etc/crontabs/root
0 4 * * * fstrim -a
45 23 * * 0-4 poweroff
